# Demo Server/my_echo.conf
start on vagrant-mounted
stop on runlevel [!2345]
exec /usr/bin/env python3.5 /vagrant/echo_server.py